# supervisord.conf
[program:php5-fpm]
command=/usr/sbin/php5-fpm -c /etc/php5/fpm

[program:mysqld]
command=/usr/bin/mysqld_safe

[program:nginx]
command=/usr/sbin/nginx
